# bootstrap/external/garagehq/ansible/templates/garage-webui.service.j2
[Unit]
Description=Garage Web UI
After=network.target

[Service]
Environment="PORT={{ garage_webui_port }}"
Environment="CONFIG_PATH=/etc/garage.toml"
ExecStart=/usr/local/bin/garage-webui
Restart=always

[Install]
WantedBy=default.target
